# trunk/pdaXrom-ng/generic/etc/udev/rules.d/95-bt-sixaxismouse.rules
SUBSYSTEM=="bluetooth" ACTION=="add" RUN+="/etc/init.d/sixaxismouse start"
SUBSYSTEM=="bluetooth" ACTION=="remove" RUN+="/etc/init.d/sixaxismouse stop"
