# tags/ps3boot-0.1/utils/99-ps3boot.rules
# tell petitboot when we see new block devices ...
SUBSYSTEM=="block",RUN+="/usr/sbin/ps3boot-udev"
