# branches/pdaXrom-ng/generic/etc/inittab
# 
# /etc/inittab
#
# $Id: inittab 1681 2004-09-01 18:12:49Z  $
#

::sysinit:/etc/init.d/rcS
tty1::respawn:/sbin/getty -L 115200 /dev/tty1 linux

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
