# trunk/minimount/95-minimount.rules
SUBSYSTEM=="block",RUN+="/usr/sbin/minimount-udev"
