# trunk/minimount/95-minimount.rules
SUBSYSTEM=="block",RUN+="/usr/lib/minimount/minimount-udev"
