# generic/etc/udev/rules.d/85-net.rules
# network devices
SUBSYSTEM=="net", ACTION=="add",    RUN+="/etc/udev/scripts/net.sh"
SUBSYSTEM=="net", ACTION=="remove", RUN+="/etc/udev/scripts/net.sh"
